# .config/systemd/user/win-scroll-volume.service
[Unit]
Description=Alt+Scroll Volume Control

[Service]
Type=simple
ExecStart=%h/.local/bin/win-scroll-volume
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
